# supervisor_conf/supervisor_mysql.conf
[program:mysqld]
command=/helpers/mysqld_supervisor.sh
autostart=true
autorestart=false
