# supervisor_conf/supervisor_apache.conf
[program:apache2]
command=/helpers/apache2_supervisor.sh
numprocs=1
autostart=true
autorestart=true